# offline-translator/com.local.hy-mt-llama-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.hy-mt-llama-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/llama-server</string>
        <string>-hf</string>
        <string>tencent/HY-MT1.5-1.8B-GGUF:Q8_0</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8110</string>
        <string>-c</string>
        <string>2048</string>
        <string>-np</string>
        <string>1</string>
        <string>--ctx-checkpoints</string>
        <string>0</string>
        <string>-ngl</string>
        <string>99</string>
    </array>

    <key>UserName</key>
    <string>__USER__</string>
    <key>GroupName</key>
    <string>staff</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/hy-mt-llama-server.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/hy-mt-llama-server.err.log</string>
</dict>
</plist>
