# offline-translator/com.local.offline-translator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.offline-translator</string>

    <!-- MLX server (Apple Silicon, Metal GPU). 跑在用户 GUI 会话里才能拿到 Metal,
         所以是 LaunchAgent 而非 system LaunchDaemon。 -->
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/share/offline-translator/venv/bin/python</string>
        <string>-m</string>
        <string>mlx_lm</string>
        <string>server</string>
        <string>--model</string>
        <string>Hy-MT2-1.8B</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8110</string>
        <string>--log-level</string>
        <string>WARNING</string>
        <!-- 必须显式给 --prompt-cache-bytes：不给的话 server.py 的裁剪分支
             (if cli_args.prompt_cache_bytes is not None) 整个不进，KV cache 的
             max_bytes 就停在默认 1<<63 = 无上限。--prompt-cache-size 只限「最多 10 个
             cache」，不限每个多大。本模型 max_position_embeddings=262144，KV 64KB/token，
             单个 cache 跑满 256k 上下文就是 16GB(实测进程涨到 17GB、峰值 20GB)。

             512M ≈ 8k tokens 的 KV。翻译是每段独立文本，共享前缀顶多一个短
             system prompt，prompt cache 复用率≈0——实测 2G vs 256M：稳态
             1891MB vs 1675MB，延迟 2.01s vs 1.77s(不降反略快)，即缩小无代价。
             真实负载 cache 只用到 ~300MB，这个上限平时根本碰不到；它防的是尾部
             ——有人翻超长文档时单个 KV cache 失控。512M 留了点余地。 -->
        <string>--prompt-cache-bytes</string>
        <string>512M</string>
    </array>

    <!-- 中性 model id "Hy-MT2-1.8B"：实为软链，指向当前量化目录(现 Hy-MT2-1.8B-mlx-4bit)。
         相对 --model + WorkingDirectory 让注册 id 干净、跨机器通用；与量化解耦——
         以后换 4bit/8bit 只改软链指向，客户端(Bob/扩展)的 model 字段不用动。 -->
    <key>WorkingDirectory</key>
    <string>__HOME__/.local/share/offline-translator</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/offline-translator.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/offline-translator.err.log</string>
</dict>
</plist>
